# meta/recipes-connectivity/bind/bind-9.9.3/bind.service
[Unit]
Description=bind
After=multi-user.target

[Service]
ExecStart=/usr/sbin/named -f
ExecReload=/usr/sbin/rndc reload
ExecStop=/usr/sbin/rndc stop

[Install]
WantedBy=multi-user.target
