# recipes-extended/logrotate/logrotate/logrotate.timer
[Unit]
Description=logrotate

[Install]
WantedBy=multi-user.target

[Timer]
OnCalendar=hourly
