# recipes-extended/logrotate/logrotate/logrotate.service
[Unit]
Description=logrotate

[Service]
Type=oneshot
ExecStart=/usr/bin/logrotate /etc/logrotate.conf
