# meta-tplino/recipes-benchmark/iperf/iperf3-3.0.3/iperf3.service
[Unit]
Description=Iperf3 Server
After=network.target

[Service]
ExecStart=/usr/bin/iperf3 -sV

[Install]
WantedBy=multi-user.target
