# recipes-benchmark/iperf/iperf-2.0.5/iperf-udp.service
[Unit]
Description=Iperf UDP Server
After=network.target

[Service]
ExecStart=/usr/bin/iperf -u -s

[Install]
WantedBy=multi-user.target
