# meta/recipes-connectivity/openssh/openssh-6.2p2/sshd.service
[Unit]
Description=OpenSSH Per-Connection Daemon
Wants=sshdgenkeys.service
After=sshdgenkeys.service

[Service]
ExecStartPre=/bin/mkdir -p /run/sshd
ExecStart=-@SBINDIR@/sshd -D
ExecReload=@BASE_BINDIR@/kill -HUP $MAINPID
StandardError=syslog
Restart=on-failure

[Install]
WantedBy=multi-user.target
