# meta-tplino/recipes-connectivity/opencapwap/files/opencapwap-wtp.service
[Unit]
Description=opencapwap wtp
After=multi-user.target
Before=hostapd.service

[Service]
ExecStart=/usr/bin/WTP -f
Type=simple
Environment="HOME=/home/root"
WorkingDirectory=/home/root
Restart=on-failure

[Install]
WantedBy=multi-user.target
