# files/mix_class/app-template.service
[Unit]
Description=@@DESCRIPTION@@

[Service]

Environment="HOME=@@APP_HOME@@"
Environment="USER_CONFIG_FILE=@@USER_CONFIG_FILE@@"
Environment="MUTABLE_DIR=@@MUTABLE_DIR@@"
ExecStart=@@APP_PREFIX@@/@@APPNAME@@/@@VERSION@@/bin/rc foreground
ExecStop=@@APP_PREFIX@@/@@APPNAME@@/@@VERSION@@/bin/rc stop
Restart=always

[Install]
WantedBy=multi-user.target
