# meta-tplino/recipes-extended/lcdproc/lcdproc/LCDd.service
[Unit]
Description=LCDd Service
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/sbin/LCDd -c /etc/lcdproc/LCDd.conf
ControlGroup=cpu:/

[Install]
WantedBy=multi-user.target
