# recipes-mgmt/nrpe/files/nrpe.socket
[Unit]
Description=NRPE Socket

[Socket]
ListenStream=5666
Accept=yes

[Install]
WantedBy=sockets.target
